# deploy/systemd/northnarrow-agent.service
[Unit]
Description=NorthNarrow XDR agent
Documentation=https://github.com/northnarrow/northnarrow-new
# Network online so the agent's outbound posture-engaged behaviour
# has a fully-up interface.
#
# NOTE: we deliberately do NOT depend on `sys-fs-bpf.mount`. /sys/fs/bpf
# is an API virtual filesystem — systemd mounts bpffs itself early in
# boot and REFUSES to materialise a sys-fs-bpf.mount unit ("Cannot
# create mount unit for API file system /sys/fs/bpf. Refusing."), even
# with an /etc/fstab entry. So `Requires=sys-fs-bpf.mount` and
# `After=…sys-fs-bpf.mount` are permanently unsatisfiable and abort the
# start job with "Unit sys-fs-bpf.mount not found". bpffs is always
# present by the time this unit runs; the ConditionPathIsMountPoint
# guard below is the real safety check.
After=network-online.target
Wants=network-online.target
# ConditionPathIsMountPoint is a soft pre-flight — units skip
# their ExecStart if the condition is unmet rather than failing.
# Pairs with the agent's own warn-and-continue when bpffs isn't
# usable, so an operator misconfig produces a clear systemd
# "skipped" status instead of a startup crash.
ConditionPathIsMountPoint=/sys/fs/bpf

[Service]
Type=notify
# Logging footprint (product guardrail). Route this unit's stdout/stderr
# into NorthNarrow's dedicated, hard-capped journal namespace instead of
# the host default journal. Stock Ubuntu's journald has ForwardToSyslog=yes
# → rsyslog → /var/log/syslog, which has NO size cap (weekly logrotate
# only): that is how a noisy bug filled 76G twice in a day. A non-default
# namespace runs its own systemd-journald@northnarrow (config:
# /etc/systemd/journald@northnarrow.conf, SystemMaxUse=1G) and by systemd
# design does NOT forward to syslog — closing the size-cap AND the
# amplification vectors with one knob, zero customer config. Read NN logs:
#   journalctl --namespace=northnarrow -u northnarrow-agent
LogNamespace=northnarrow
ExecStart=/usr/local/bin/northnarrow-agent \
    --combat-rules /etc/northnarrow/combat-rules.v4 \
    --admin-pub /etc/northnarrow/admin.pub \
    --admin-socket /run/northnarrow/admin.sock \
    --pid-file /run/northnarrow/agent.pid
# Watchdog (Tappa 7 task 6 W4) owns restart policy. systemd MUST
# NOT auto-restart the agent — doing so would race the
# watchdog's restart-backoff loop and double-spawn agents on
# every death. Per design §2.3 "restart policy lives only in
# the watchdog — the agent's systemd unit uses Restart=no and
# delegates respawn to the watchdog."
Restart=no
RuntimeDirectory=northnarrow
RuntimeDirectoryMode=0700

# Resource ceilings (Beta Step 2a — confine an ADE runaway to the
# agent's own cgroup instead of letting the kernel OOM-killer roam
# the host and reap sshd / other critical services — "Family A"
# resource-exhaustion cascade).
#
#   MemoryHigh — soft throttle. Above this the cgroup is put under
#     reclaim pressure (and briefly stalled) but NOT killed. Set to
#     9G: the in-process ADE (Foundation-Sec-8B Q4_K_M) peaks ~7 GB
#     RSS during steady-state inference (see docs/ADE_BACKEND_NOTES.md
#     "Peak RSS 7.0 GB"), so 9G leaves the legitimate model + agent
#     working set unthrottled while applying back-pressure the moment
#     usage climbs past it.
#   MemoryMax — hard ceiling. The cgroup's own OOM-killer fires here,
#     killing ONLY the agent (the watchdog then respawns it per its
#     own restart policy). Set to 11G: ~1.5× the legitimate peak for
#     transient spikes, and on the 15.24 GiB reference host (Hetzner
#     CCX23) it still leaves ~4 GiB for the OS, sshd and the watchdog
#     — so a runaway can never starve the host into a cascade.
#     On a smaller VM the agent's pre-ADE RAM guard (Beta Step 2b)
#     skips model load entirely and runs detection-only, so this
#     ceiling is never approached there.
#   TasksMax — fork/thread-leak backstop. The agent runs a tokio
#     runtime (a few workers), candle's CPU kernels (≈ physical_cores
#     threads), plus FIM / net / admin tasks — a few dozen at most.
#     512 is well above any legitimate count yet caps a runaway.
#
# NOTE: deliberately NO CPUQuota. Unlike the watchdog (idle until the
# agent dies, so capped at 10%), the agent does latency-sensitive
# detection and CPU-bound ADE inference that wants every available
# core; throttling CPU would slow verdicts without bounding the OOM
# vector this block exists to contain. Memory is the cascade risk, so
# memory is what we cap.
MemoryHigh=9G
MemoryMax=11G
TasksMax=512

# Hardening — `ProtectSystem=strict` makes the entire filesystem
# read-only EXCEPT the ReadWritePaths whitelist below. NoNewPrivileges
# kills any latent setuid-binary escalation in case the agent ever
# spawns one.
#
# NOTE: /sys/fs/bpf/northnarrow is intentionally NOT listed. Listing a
# path under the /sys/fs/bpf API VFS makes systemd inject an implicit
# Requires=/After= on the (unsatisfiable) sys-fs-bpf.mount unit — the
# same failure mode as the explicit dep removed above. It is also
# unnecessary: ProtectSystem=strict EXEMPTS /dev, /proc and /sys from
# the read-only remount (only ProtectKernelTunables= would make /sys
# read-only, and we do not set it), so the agent retains write access
# to its bpffs pin directory without a ReadWritePaths entry.
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=read-only — deliberately NOT =yes/=true. The FIM
# credential-read rules (NN-L-FIM-011..017: AWS/Azure/GCP/Docker/browser/
# KeePass/GnuPG, plus the /root + /home SSH-key & shell-history watches)
# depend on the agent being able to stat(2) those files at boot to record
# their (dev,ino) in the kernel WATCHED_PATHS map — the BPF observe hooks
# only emit a FIM event for an inode that is IN that map. `ProtectHome=yes`
# mounts an EMPTY tmpfs over /home, /root and /run/user inside this unit's
# mount namespace, so every credential stat() returns ENOENT, those inodes
# never enter WATCHED_PATHS, and the credential-theft rules fire NOTHING —
# a silent, total coverage hole (the kernel BPF side still sees the opens;
# it simply has no watched inode to match them against). `read-only`
# restores visibility while keeping the write-hardening intact: the agent
# only ever stat()s/reads these paths — it NEVER writes under a home dir
# (all writes go to the ReadWritePaths below) — so read-only costs us
# nothing and is verifiable with `systemd-analyze security`. Do NOT revert
# to =yes without first removing every /root + /home entry from
# configs/fim-paths.v1; the boot-time regression guard in
# agent/src/fim/attach.rs (populate_watched_paths) WARNs if this is ever
# masked again.
ProtectHome=read-only
ReadWritePaths=/run/northnarrow /var/lib/northnarrow
# CAP_BPF        — load+pin BPF maps/programs (Tappa 7 task 5/6/7)
# CAP_SYS_ADMIN  — bpf(2) historically needed this; CAP_BPF
#                  separated it on 5.8+ but kernel still accepts
#                  the union, so we keep both for compat
# CAP_NET_ADMIN  — iptables-restore on COMBAT entry (Tappa 7 task 7)
# CAP_DAC_OVERRIDE — chattr +i on /var/lib/northnarrow (Tappa 7 task 5)
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target
