# deploy/systemd/northnarrow-agent.service
[Unit]
Description=NorthNarrow XDR agent
Documentation=https://github.com/northnarrow/northnarrow-new
# Network online so the agent's outbound posture-engaged behaviour
# has a fully-up interface.
#
# NOTE: we deliberately do NOT depend on `sys-fs-bpf.mount`. /sys/fs/bpf
# is an API virtual filesystem — systemd mounts bpffs itself early in
# boot and REFUSES to materialise a sys-fs-bpf.mount unit ("Cannot
# create mount unit for API file system /sys/fs/bpf. Refusing."), even
# with an /etc/fstab entry. So `Requires=sys-fs-bpf.mount` and
# `After=…sys-fs-bpf.mount` are permanently unsatisfiable and abort the
# start job with "Unit sys-fs-bpf.mount not found". bpffs is always
# present by the time this unit runs; the ConditionPathIsMountPoint
# guard below is the real safety check.
After=network-online.target
Wants=network-online.target
# ConditionPathIsMountPoint is a soft pre-flight — units skip
# their ExecStart if the condition is unmet rather than failing.
# Pairs with the agent's own warn-and-continue when bpffs isn't
# usable, so an operator misconfig produces a clear systemd
# "skipped" status instead of a startup crash.
ConditionPathIsMountPoint=/sys/fs/bpf

[Service]
Type=notify
ExecStart=/usr/local/bin/northnarrow-agent \
    --combat-rules /etc/northnarrow/combat-rules.v4 \
    --admin-pub /etc/northnarrow/admin.pub \
    --admin-socket /run/northnarrow/admin.sock \
    --pid-file /run/northnarrow/agent.pid
# Watchdog (Tappa 7 task 6 W4) owns restart policy. systemd MUST
# NOT auto-restart the agent — doing so would race the
# watchdog's restart-backoff loop and double-spawn agents on
# every death. Per design §2.3 "restart policy lives only in
# the watchdog — the agent's systemd unit uses Restart=no and
# delegates respawn to the watchdog."
Restart=no
RuntimeDirectory=northnarrow
RuntimeDirectoryMode=0700

# Hardening — `ProtectSystem=strict` makes the entire filesystem
# read-only EXCEPT the ReadWritePaths whitelist below. NoNewPrivileges
# kills any latent setuid-binary escalation in case the agent ever
# spawns one.
#
# NOTE: /sys/fs/bpf/northnarrow is intentionally NOT listed. Listing a
# path under the /sys/fs/bpf API VFS makes systemd inject an implicit
# Requires=/After= on the (unsatisfiable) sys-fs-bpf.mount unit — the
# same failure mode as the explicit dep removed above. It is also
# unnecessary: ProtectSystem=strict EXEMPTS /dev, /proc and /sys from
# the read-only remount (only ProtectKernelTunables= would make /sys
# read-only, and we do not set it), so the agent retains write access
# to its bpffs pin directory without a ReadWritePaths entry.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/northnarrow /var/lib/northnarrow
# CAP_BPF        — load+pin BPF maps/programs (Tappa 7 task 5/6/7)
# CAP_SYS_ADMIN  — bpf(2) historically needed this; CAP_BPF
#                  separated it on 5.8+ but kernel still accepts
#                  the union, so we keep both for compat
# CAP_NET_ADMIN  — iptables-restore on COMBAT entry (Tappa 7 task 7)
# CAP_DAC_OVERRIDE — chattr +i on /var/lib/northnarrow (Tappa 7 task 5)
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target
